# launchagent/com.session-indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.session-indexer</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python installation -->
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>session_index.indexer</string>
        <string>--incremental</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/session-indexer.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/session-indexer.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
